# hive/sonm-standby.service
[Unit]
Description=Switcher for Sonm/Hive internal mining
After=sonm-worker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/sonm-standby
